# systemd/zfs-scrub.service
[Unit]
Description=Performs ZFS pools scrub

[Service]
Type=simple
ExecStart=/sbin/zpool scrub system
