# old-dotfiles/systemd/zfs-monitor.timer
[Unit]
Description=Send ZFS status email weekly

[Timer]
OnBootSec=1min
OnUnitActiveSec=1w

[Install]
WantedBy=timers.target
